# linux-os/seccam.cron
#
10 22 * * 4 root /etc/motion/c-allstart
15 15 * * 4 root /etc/motion/c-allstop
20 17 * * 4 root /etc/motion/c-allstart
10 20 * * 4 root /etc/motion/c-allstop
#
10 02 * * 7 root /etc/motion/c-allstart
15 15 * * 7 root /etc/motion/c-allstop
20 17 * * 7 root /etc/motion/c-allstart
10 20 * * 7 root /etc/motion/c-allstop
#
