# linux-os/etc/cron.d/seccam
# NVR napi videók mentése

0 0 * * *      root /usr/local/bin/seccam 2>>/dev/null >>/dev/null
*/5 * * * *   root /usr/local/bin/seccam-serv 2>>/dev/null >>/dev/null

# log váltása

0 0 * * *      root printf "\n" >>/var/log/motion/cam.log
0 0 * * 0      root /etc/motion/c-logrotate

#
